# com.stockview.update-yahoo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- タスク名 -->
    <key>Label</key>
    <string>com.stockview.update-yahoo</string>

    <!-- 実行するコマンド -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ogurahirokazu/claude-code-lab/stock-view/update-screening-data.py</string>
        <string>--mode</string>
        <string>yahoo</string>
    </array>

    <!-- 作業ディレクトリ -->
    <key>WorkingDirectory</key>
    <string>/Users/ogurahirokazu/claude-code-lab/stock-view</string>

    <!-- 毎週土曜日 5:00 に実行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>6</integer>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- ログ出力先 -->
    <key>StandardOutPath</key>
    <string>/Users/ogurahirokazu/claude-code-lab/stock-view/logs/update-yahoo.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ogurahirokazu/claude-code-lab/stock-view/logs/update-yahoo-error.log</string>

    <!-- 環境変数（PATHを明示的に設定） -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
